# deploy/systemd/cloudflared-lizo.service.example
[Unit]
Description=Cloudflare Tunnel for Lizo
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
ExecStart=/usr/local/bin/cloudflared tunnel --config /home/pi/Lizo/deploy/cloudflared/config.yml run
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
